# root/etc/services.d/web-admin/run
#!/usr/bin/with-contenv bash

cd /app/web-admin || exit

exec \
	s6-setuidgid abc npm run start
